# synergy/synergy.service
[Unit]
Description=Synergy

[Service]
ExecStart=/usr/bin/synergy --no-tray

[Install]
WantedBy=default.target
